# hhvm/hhvm.sv.conf
[program:hhvm]
command=hhvm --mode daemon -vServer.Type=fastcgi -vServer.Port=9000
redirect_stderr=true
