# nginx/nginx.sv.conf
[program:nginx]
command=
redirect_stderr=true
