# mariadb/mariadb.sv.conf
[program:mariadb]
command=/usr/bin/mysqld_safe
redirect_stderr=true
